# green/conf/init_script.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/start.sh

[Install]
WantedBy=multi-user.target